# rootfs/etc/cron.d/crontabs/root
10 0-23 * * * /usr/local/bin/fetchWeather.sh
