# host-rootfs/etc/udev/rules.d/97-parrot-df3120.rules
ACTION=="add", ATTRS{idVendor}=="1d6b", ATTRS{idProduct}=="0102", NAME="usb%n", RUN+="/etc/udev/scripts/parrot-df3120.sh usb%n"
