# rootfs/var/spool/cron/crontabs/root
10 */1 * * * /usr/local/bin/fetchWeather.sh
@reboot /usr/local/bin/fetchWeather.sh
